# crontab/mac/crontab
*/30 * * * * $HOME/scripts/caffeinate.sh >> $HOME/scripts/log/caffeinate.log 2>&1
