# crontab/mac/crontab
*/30 * * * * /Users/i863328/scripts/caffeinate.sh >> /Users/i863328/scripts/log/caffeinate.log 2>&1
